# etc/systemd/system/mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.service
# mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.service
# autogenerated by Podman 2.0.4
# Sat Aug 22 20:44:02 HKT 2020

[Unit]
Description=Podman mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
BindsTo=mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.service
After=mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/bin/rm -f %t/mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.pid %t/mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.ctr-id
ExecStart=/usr/bin/podman run --conmon-pidfile %t/mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.pid --cidfile %t/mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.ctr-id --cgroups=no-conmon --pod-id-file %t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pod-id -d --replace --log-driver journald --name mydashboard-loki --rm --user 65534:65534 --volume /etc/mydashboard/loki:/etc/loki --volume /var/lib/mydashboard/loki:/var/lib/loki docker.io/grafana/loki --config.file /etc/loki/loki.yaml
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.ctr-id
PIDFile=%t/mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
