# etc/systemd/system/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.service
# mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.service
# autogenerated by Podman 2.0.4
# Sat Aug 22 20:44:02 HKT 2020

[Unit]
Description=Podman mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
Requires=mdd-c-2d6b71c02c52ed6984aeede9c77f41e984cbe1185c6d002b501a484a642fd4fa.service mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.service mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.service mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.service
Before=mdd-c-2d6b71c02c52ed6984aeede9c77f41e984cbe1185c6d002b501a484a642fd4fa.service mdd-c-307e49b00b1e69db995749c19e78a0dc75a3ecc87528a37598325d92c1092a2f.service mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.service mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/bin/rm -f %t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pid %t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pod-id
ExecStartPre=/usr/bin/podman pod create --infra-conmon-pidfile %t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pid --pod-id-file %t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pod-id --name mydashboard --network mydashboard --publish 127.0.0.1:3100:3100 --replace
ExecStart=/usr/bin/podman pod start --pod-id-file %t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pod-id
ExecStop=/usr/bin/podman pod stop --ignore --pod-id-file %t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pod-id -t 10
ExecStopPost=/usr/bin/podman pod rm --ignore -f --pod-id-file %t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pod-id
PIDFile=%t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pid
KillMode=none
Type=forking

[Install]
Alias=mdd.service
WantedBy=multi-user.target default.target
