# etc/systemd/system/mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.service
# mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.service
# autogenerated by Podman 2.0.4
# Sat Aug 22 20:44:02 HKT 2020

[Unit]
Description=Podman mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
BindsTo=mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.service
After=mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/bin/rm -f %t/mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.pid %t/mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.ctr-id
ExecStartPre=/bin/rm -f /run/mydashboard/grafana/grafana.sock
ExecStart=/usr/bin/podman run --conmon-pidfile %t/mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.pid --cidfile %t/mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.ctr-id --cgroups=no-conmon --pod-id-file %t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pod-id -d --replace --log-driver journald --name mydashboard-grafana --rm --volume /etc/mydashboard/grafana:/etc/grafana --volume /var/lib/mydashboard/grafana:/var/lib/grafana --volume /run/mydashboard/grafana:/run/grafana docker.io/grafana/grafana
ExecStartPost=/bin/sh -c "i=0 & until [[ -S /run/mydashboard/grafana/grafana.sock ]]; do i=$(($i+1)); sleep 0.25; [[ $i == '200' ]] && exit 1; done; chmod 0666 /run/mydashboard/grafana/grafana.sock"
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.ctr-id
PIDFile=%t/mdd-c-6595a2227bf3dcbd70db51cdd5a9497fc2ac19495ab2248a0179f4cfe47bee61.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
