# etc/systemd/system/mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.service
# mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.service
# autogenerated by Podman 2.0.4
# Sat Aug 22 20:44:02 HKT 2020

[Unit]
Description=Podman mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
BindsTo=mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.service
After=mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
ExecStartPre=/bin/rm -f %t/mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.pid %t/mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.ctr-id
ExecStart=/usr/bin/podman run --conmon-pidfile %t/mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.pid --cidfile %t/mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.ctr-id --cgroups=no-conmon --pod-id-file %t/mdd-p-ceb67724093c6e374e995da00d5a10ac90190ef084e3f1e9ca03ab0b2041ebd1.pod-id -d --replace --log-driver journald --name mydashboard-prometheus --rm --volume /etc/mydashboard/prometheus:/etc/prometheus --volume /var/lib/mydashboard/prometheus/data:/prometheus docker.io/prom/prometheus
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.ctr-id
PIDFile=%t/mdd-c-f2059bf1c3b6352713266572069147ce1c2f80e252c0293453aaf21f569ab1e5.pid
KillMode=none
Type=forking

[Install]
WantedBy=multi-user.target default.target
